# com.newsvibe.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.newsvibe.digest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <!-- Замени __PROJECT_DIR__ на абсолютный путь к корню проекта.
         Проще через sed при установке — см. README, секция launchd. -->
    <string>__PROJECT_DIR__/run_digest.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd даёт минимальный PATH; добавляем места claude/yt-dlp/node.
         Скорректировать под факт: `which claude yt-dlp node`. -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/com.newsvibe.digest.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.newsvibe.digest.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
